# deploy/phase-server.service
[Unit]
Description=phase-server: multiplayer game server for phase.rs
After=network.target

[Service]
Type=simple
User=phase
Group=phase
WorkingDirectory=/opt/phase-server
ExecStart=/opt/phase-server/phase-server --port 9374 --data-dir /opt/phase-server/data --cors-origin *
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
LogsDirectory=phase-rs
Environment=PHASE_LOG_DIR=/var/log/phase-rs
Environment=RUST_LOG=info

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/phase-server
PrivateTmp=true

[Install]
WantedBy=multi-user.target
